# systemd/APP_NAME.production.service.d/local.conf
[Service]
Environment="DOCKER_COMPOSE_FILE="
Environment="FLASK_DB_DATABASE="
Environment="FLASK_DB_PASSWORD="
Environment="FLASK_SECRET_KEY="
Environment="SSL_LOCATION="
