# systemd/APP_NAME.production.service.d/local.conf
[Service]
Environment="DOCKER_COMPOSE_FILE="
Environment="FLASK_SECRET_KEY="
Environment="SSL_LOCATION="
Environment="SOLUTIONS_FOLDER="
